# usr/lib/systemd/system/godoc-runner.socket
[Socket]
BindToDevice=lo
ListenStream=6062
Accept=yes

[Install]
WantedBy=sockets.target
