# usr/lib/systemd/system/godoc-builder.socket
[Socket]
BindToDevice=lo
ListenStream=6061

[Install]
WantedBy=sockets.target
